# rootfs/etc/rc.local
#!/bin/sh
# /etc/rc.local - Local system initialization script

# Set up environment
export PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Create necessary directories
mkdir -p /var/log
mkdir -p /var/run
mkdir -p /var/tmp

# Set permissions
chmod 1777 /tmp
chmod 1777 /var/tmp

# System initialization tasks
echo "System initialized at $(date)" > /var/log/boot.log
